# scripts/content-engine/trigger/com.redrabbit.mediachecker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd LaunchAgent: prueft alle 30 Minuten ob ein genehmigter Artikel
     auf seinen Media-Schritt wartet. Triggert maximal einmal pro Tag.
     Install:
       cp scripts/content-engine/trigger/com.redrabbit.mediachecker.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.redrabbit.mediachecker.plist -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.redrabbit.mediachecker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>$HOME/dev/redrabbit-daily/scripts/content-engine/trigger/run-media-check.sh</string>
    </array>

    <!-- Alle 30 Minuten (1800 Sekunden). -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Sofort beim Laden/Login starten (z.B. nachdem der Mac eingeschaltet wurde, waehrend eine
         Freigabe-Mail offline bestaetigt wurde) UND danach alle 30 Min. Der Lauf ist idempotent
         (Tages-Stempel + reagiert nur auf den heutigen Marker), also unschaedlich wenn nichts ansteht. -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/redrabbit-mediachecker.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/redrabbit-mediachecker.err</string>
</dict>
</plist>
